# roles/pinger/templates/pinger.service.j2
[Unit]
Description=Pinger service
After=network-online.target

[Service]
User=ubuntu
ExecStart=/usr/local/bin/pinger

[Install]
WantedBy=multi-user.target
