# roles/pinger/templates/pinger.service.j2
[Unit]
Description=Pinger service

[Service]
User=pinger
ExecStart=/usr/local/bin/pinger

[Install]
WantedBy=multi-user.target
